# config/apache2.conf
[program:apache2]
command=/usr/local/bin/apache2-foreground
autostart=true
autorestart=true
startretries=1
startsecs=1
